# .claude/launchd/com.wally.hermes-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.wally.hermes-daemon.plist
  Keeps the Hermes Telegram bot alive on macOS login.

  ProgramArguments invokes adapters/hermes/hermes_daemon_wrapper.sh, which:
    - Auto-detects the hermes binary (Homebrew /opt/homebrew, /usr/local, etc.)
    - Auto-detects the daemon subcommand (serve|bot|daemon|run)
    - Sources $REPO/.env for secrets (NOTION_API_KEY, etc.)
    - Rotates logs/hermes-daemon.log when it grows past 10 MB
  This means the plist itself never needs editing across machines.

  Install:    make hermes-daemon-install
  Uninstall:  make hermes-daemon-uninstall
  Restart:    make hermes-restart
  Logs:       make hermes-logs
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.wally.hermes-daemon</string>

    <!-- Wrapper handles binary + subcommand detection. No more hardcoded paths. -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/josecampos/Documents/wally-trader/adapters/hermes/hermes_daemon_wrapper.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/josecampos/Documents/wally-trader</string>

    <!-- Start automatically on user login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart if the process crashes or exits non-zero -->
    <key>KeepAlive</key>
    <true/>

    <!-- Throttle rapid restart loops (seconds between restart attempts) -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <!-- Logs inside the repo so make hermes-logs can easily tail them -->
    <key>StandardOutPath</key>
    <string>/Users/josecampos/Documents/wally-trader/logs/hermes-daemon.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/josecampos/Documents/wally-trader/logs/hermes-daemon.log</string>

    <!-- Forward PATH so the wrapper can find hermes / node / python -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/josecampos</string>
    </dict>

    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
